# config/supervisord/up/0/up196.ini
[program:up196]
command=/root/.virtualenvs/poke/bin/python /var/poke/pokemap/manage.py -c ./config/production.py up -p 196
user=root
autorestart=true
stdout_logfile=/var/log/supervisor/up196-supervisord.log ; 標準出力ログ
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
redirect_stderr=true
